# .tests/borealis/etc/systemd/system/maitred.service
[Unit]
Description=maitred

[Service]
ExecStart=/opt/google/cros-containers/bin/maitred

[Install]
WantedBy=multi-user.target
