# .tests/borealis/etc/systemd/system/vshd.service
[Unit]
Description=vshd

[Service]
ExecStart=/opt/google/cros-containers/bin/vshd

[Install]
WantedBy=basic.target
